# ops/cron/theme-browser-registry-refresh.cron
0 3 * * 0 cd /opt/theme-browser-monorepo && /usr/bin/env bash ./scripts/registry-refresh.sh >> /var/log/theme-browser-registry-refresh.log 2>&1
