# ops/systemd/theme-browser-registry-refresh.service
[Unit]
Description=Refresh theme-browser registry data
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
EnvironmentFile=/etc/theme-browser/registry-refresh.env
ExecStart=/usr/bin/docker run --rm --name theme-browser-registry-refresh --env-file /etc/theme-browser/registry-refresh.env -v /var/lib/theme-browser-refresh:/var/lib/theme-browser-refresh theme-browser-registry-refresh:local

[Install]
WantedBy=multi-user.target
