# ops/systemd/theme-browser-registry-refresh.service
[Unit]
Description=Refresh theme-browser registry data
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/opt/theme-browser-monorepo
EnvironmentFile=/etc/theme-browser/registry-refresh.env
ExecStart=/usr/bin/env bash ./scripts/registry-refresh.sh

[Install]
WantedBy=multi-user.target
